# BusyBox/MyConfig/etc/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init
::shutdown:/sbin/swapoff -a
::shutdown:sync
# ::shutdown:umount -a -r
::ctrlaltdel:/sbin/reboot
::respawn:-/bin/sh
tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh
tty4::askfirst:-/bin/sh

