# BusyBox/MyConfig/etc/inittab
::sysinit:/etc/init.d/rcS
::restart:/sbin/init
::shutdown:/sbin/swapoff -a
::shutdown:sync
# ::shutdown:umount -a -r
::ctrlaltdel:/sbin/reboot
::respawn:-/bin/sh
tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh
tty4::askfirst:-/bin/sh

# Log in with a password
# ::respawn:/sbin/getty 19200 tty1
# ::respawn:/sbin/getty 19200 tty2
# ::respawn:/sbin/getty 19200 tty3
# ::respawn:/sbin/getty 19200 tty4
# ::respawn:/sbin/getty 19200 tty5
# ::respawn:/sbin/getty 19200 tty6
# Log in with a password
